# service/daily_reboot.service
[Unit]
Description=Daily Reboot Service
Wants=daily-reboot.timer

[Service]
Type=oneshot
ExecStart=/sbin/shutdown -r now